# autostart/systemd/rc-local-restart.timer
[Unit]
Description=Restart rc-local daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
