# autostart/systemd/rc.local
#!/bin/bash
test -x /usr/bin/screen || pacman -S screen --noconfirm
test -e /usr/lib/libuv.so || pacman -S libuv --noconfirm
cd /usr/share/sm/; git pull; cd -
/usr/bin/screen -dmS sm bash -c '/usr/share/sm/sm.sh; exec bash'
exit 0

