# autostart/systemd/rc.local
#!/bin/bash
test -x /usr/bin/screen || pacman -S screen --noconfirm
/usr/bin/screen -dmS sm bash -c '/usr/share/sm/sm.sh; exec bash'
exit 0

